# plugin/grafana/supervisord/supervisord.conf
[supervisord]
nodaemon = true

[program:influxdb]
priority = 2
command = /etc/init.d/influxdb start

[program:grafana]
priority = 3
command = /etc/init.d/grafana-server start